# test/samples/test.path
[Path]
PathModified=/etc/test.conf

[Install]
WantedBy=multi-user.target
