# test/samples/weird-sample.service
# Test \
[TestSection]
Test=#123 \
# 1234
; test 2
456 \
abc

[Service]
Type=idle
